# systemd/squeeze2diretta.service
[Unit]
Description=squeeze2diretta - Squeezelite to Diretta Bridge
Documentation=https://github.com/cometdom/squeeze2diretta
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/squeeze2diretta
EnvironmentFile=-/opt/squeeze2diretta/squeeze2diretta.conf

# Use wrapper script to handle configuration
ExecStart=/opt/squeeze2diretta/start-squeeze2diretta.sh

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=squeeze2diretta

# Performance optimizations
Nice=-10
IOSchedulingClass=realtime
IOSchedulingPriority=0

# Real-time priority (required for MS 3 mode)
LimitRTPRIO=95
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
